# deployment/systemd/stockstory-backend.service
[Unit]
Description=StockStory Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/PREDICTION-ENGINE
ExecStart=/usr/bin/npm run prod:backend
Restart=always
RestartSec=5
Environment="NODE_ENV=production"
Environment="PORT=4001"
Environment="HOST=0.0.0.0"
Environment="SGLANG_INTELLIGENCE_URL=http://localhost:11434"

[Install]
WantedBy=multi-user.target
